# bin/neopixel.service
[Unit]
Description=Service for Neopixel driver.
# After=multi-user.target
After=multi-user.target redis-server.service

[Service]
User=pi
Type=idle
ExecStart=sudo /usr/bin/python3 -u neopixel_startup.py
WorkingDirectory=/home/pi/new_infrastructure/stream-sim-backend/bin
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
