# bin/streamsim.service
[Unit]
Description=Service for StreamSim.
# After=multi-user.target
After=multi-user.target network.target redis-server.service

[Service]
User=pi
Type=idle
ExecStart=sudo /usr/bin/python3 -u main.py tektrain_real
WorkingDirectory=/home/pi/new_infrastructure/stream-sim-backend/bin
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
