# bin/streamsim.service
[Unit]
Description=Service for StreamSim.
# After=multi-user.target
After=multi-user.target network.target redis-server.service

[Service]
User=pi
Type=idle
ExecStart=/usr/bin/python3 -u main.py elsa
WorkingDirectory=/home/pi/v2/stream-sim-backend/stream_simulator
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
